# file-editor/rootfs/etc/services.d/file-editor/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Closing File Editor"

rm -rf /var/run/editor.lock
