# file-editor/rootfs/etc/services.d/file-editor/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

node /app/server.js