# minecraft-server/rootfs/etc/services.d/mcpanel/run
#! /usr/bin/with-contenv bashio
# shellcheck shell=bash

cd /mcpanel || exit
bashio::log.info "Staring MCpanel"
node server.js