# minecraft-server/rootfs/etc/services.d/sshd/run
#! /usr/bin/with-contenv bashio
# shellcheck shell=bash
mkdir -p /run/sshd
/usr/sbin/sshd